# etc/supervisor/supervisor.conf
[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:server]
command=python /code/node/main.py
autorestart=true
